# upload_recordings.service
[Unit]
Description=Upload Recordings Service
After=network.target

[Service]
Environment="MYSQL_HOST=jambonz-xyz-aurora-mysql-cluster.xxxxxxx"
Environment="MYSQL_USER=xxxx"
Environment="MYSQL_PASSWORD=xxxxx"
Environment="MYSQL_DATABASE=jambones"
Environment="BASIC_AUTH_USERNAME=xxxxx"
Environment="BASIC_AUTH_PASSWORD=xxxxx"
Environment="ENCRYPTION_SECRET=xxxxxx"
# When installed via the .deb, the binary lives at /usr/bin/upload_recordings
# and its bundled .so dependencies are at /usr/lib/upload-recordings/, found
# automatically via /etc/ld.so.conf.d/upload-recordings.conf. No LD_LIBRARY_PATH
# is required.
ExecStart=/usr/bin/upload_recordings --port 3017
Restart=always
RestartSec=5
User=root
Group=daemon
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
